# assignment04/99-hello_keyboard.rules
# /etc/udev/rules.d/99-hello_keyboard.rules

#ACTION=="add", SUBSYSTEM=="usb", ENV{ID_INPUT_KEYBOARD}=="1", RUN+="/usr/sbin/insmod /home/dpoveda/git/little-penguin-1/assignment04/main.ko"

ACTION=="add", SUBSYSTEM=="usb", ATTR{bInterfaceClass}=="03", ATTR{bInterfaceProtocol}=="01", RUN+="/usr/sbin/insmod /home/dpoveda/git/little-penguin-1/assignment04/main.ko"
